# backend/conf/supervisor/supervisord.d/myapp.ini
[program:my_program]
command=/root/.local/share/virtualenvs/backend-2ZtQ7vfC/bin/uwsgi --ini /usr/local/share/helloBlog/backend/conf/uwsgi.ini
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/supervisor.log
